# src/scripts/northern-hemisphere-prod.service
[Unit]
Description=Northern Hemisphere Service
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/northern-hemisphere-prod
Restart=on-abort

[Install]
WantedBy=multi-user.target